# supervisor/caddy.conf
[program:caddy]
command=/usr/local/bin/caddy -agree
directory=/home/skylines/src/
user=skylines
stdout_logfile=/home/skylines/caddy/out.log
stderr_logfile=/home/skylines/caddy/err.log
